# ROMFS/RC/px4fmu-v2_follow_fc/init.d/rc.local
#!nsh
#
# local startup script
#

mavlink start -d /dev/ttyS1 -b 19200

if waypoint_fc start
then
fi

