# ROMFS/RC/px4fmu-v2_follow_wb/init.d/rc.local
#!nsh
#
# local startup script
#

if waypoint_wb start
then
fi

