# ROMFS/RC/px4fmu-v2_follow_wb/init.d/rc.local
#!nsh
#
# local startup script
#

mavlink start -d /dev/ttyS2 -b 19200

if waypoint_wb start
then
fi

