# ROMFS/extras/rc.local.sample
#!nsh
#
# local startup script
#

#mavlink start -d /dev/ttyS2 -b 57600
